# roles/librenms_app/templates/librenms-ha-startup-repair.service.j2
[Unit]
Description=LibreNMS HA startup repair
Wants=network-online.target remote-fs.target
After=network-online.target remote-fs.target
After={{ librenms_glusterd_service_name }}.service
After={{ librenms_mariadb_service_name }}.service
After={{ librenms_haproxy_service_name }}.service
After={{ librenms_redis_service_name }}.service
After={{ librenms_redis_sentinel_service_name_effective }}.service

[Service]
Type=oneshot
ExecStart={{ librenms_startup_repair_path }}

[Install]
WantedBy=multi-user.target
